# deploy/chromehearts-monitor.timer
[Unit]
Description=Run Chrome Hearts appointment check every minute

[Timer]
OnCalendar=minutely
# If a check is still running when the next tick fires, the tick is a
# no-op (systemd won't start an already-active unit), so runs never pile up.
AccuracySec=1s

[Install]
WantedBy=timers.target
