# deploy/chromehearts-heartbeat.service
[Unit]
Description=Daily liveness ping for the Chrome Hearts monitor

[Service]
Type=oneshot
User=chmonitor
EnvironmentFile=/etc/chromehearts-monitor.env
ExecStart=/bin/sh -c 'curl -fsS -H "X-Priority: min" -H "X-Title: Appointment watchers" -d "Heartbeat from $(hostname): chromehearts=$(systemctl is-active chromehearts-monitor.service), goldendiner=$(systemctl is-active resy-monitor.service)" "https://ntfy.sh/$NTFY_TOPIC" > /dev/null'
