# deploy/chromehearts-monitor.service
[Unit]
Description=Chrome Hearts appointment watcher (continuous)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=chmonitor
WorkingDirectory=/opt/chromehearts/run
EnvironmentFile=/etc/chromehearts-monitor.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/chromehearts/venv/bin/python /opt/chromehearts/monitor.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
